# docker/s6-rc.d/pulse/run
#!/command/with-contenv sh

php /var/www/artisan pulse:check
