# docker/s6-rc.d/reverb/run
#!/command/with-contenv sh

php /var/www/artisan reverb:start
